# deploy/systemd/phantom-web.service
[Unit]
Description=Phantom Web
Documentation=https://github.com/layzndev/phantom
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=phantom
Group=phantom
WorkingDirectory=/opt/phantom/apps/web
EnvironmentFile=/etc/phantom/phantom-web.env
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30
SyslogIdentifier=phantom-web

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/phantom/apps/web/.next

[Install]
WantedBy=multi-user.target
