# deploy/phantom-agent.service
[Unit]
Description=Phantom Docker Agent
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/phantom
EnvironmentFile=/opt/phantom/apps/agent/.env
ExecStart=/usr/bin/npm run start --workspace @phantom/agent
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
